# server/ops/systemd/mesh-health.timer
[Unit]
Description=Check MeshChat relay health every 15 minutes

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=60
Persistent=true
Unit=mesh-health.service

[Install]
WantedBy=timers.target
