# server/ops/systemd/mesh-health.service
[Unit]
Description=MeshChat relay health check
After=mesh-server.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/mesh_messenger
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/mesh_messenger/.venv/bin/python /root/mesh_messenger/server/ops/healthcheck_server.py
Nice=10
UMask=0077
